# .config/systemd/user/obs_replay_notif.service
[Unit]
Description=obs_replay_notif

[Service]
ExecStart=obs_replay_notif

[Install]
WantedBy=default.target
